# systemd/gpu-telemetry-flush.service
[Unit]
Description=GPU Telemetry Flush Spool (gpu-telemetry)

[Service]
Type=oneshot
WorkingDirectory=%h/projects/gpu-telemetry
EnvironmentFile=%h/projects/gpu-telemetry/.env
ExecStart=%h/projects/gpu-telemetry/bin/uv.sh run %h/projects/gpu-telemetry/bin/flush_spool.py
